# media-sound/mpdscribble/files/mpdscribble
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

: CFGFILE=${CFGFILE:=/etc/mpdscribble.conf}

get_config() {
	x=$1
	test -e ${CFGFILE} || return 1
	cat ${CFGFILE} | \
		sed -n -e '/^[ \t]*'${x}'/{s:^[ \t]*'${x}'[ \t]\+"\?\([^#"]\+\)[^"]*"\?$:\1: ; p }' | \
		awk '{print $3}'
}

extra_started_commands='reload'
command=/usr/bin/mpdscribble
#command_args="--conf ${CFGFILE}"
required_files=${CFGFILE}
pidfile=$(get_config pid)
description="An MPD client that submits information to Audioscrobbler"

depend() {
	need mpd
	use net
	config ${CFGFILE}
}

reload() {
	ebegin "Reloading ${RC_SVCNAME}"
	start-stop-daemon --pidfile ${pidfile} --signal HUP
	eend $?
}
